# init.d/hostapd_action
#!/bin/sh /etc/rc.common
config_load hostapd_action

USE_PROCD=1
START=31
STOP=20
config_get IFACE network IFACE
PIDFILE=/var/run/hostapd_cli_$IFACE.pid

ACTION=/etc/hostapd_action

start_service() {
#	if [ -f $PIDFILE ] && kill -0 $(cat $PIDFILE); then
#		echo 'Service already running'
#		return 1
#	fi
	echo starting
	hostapd_cli -i $IFACE -a $ACTION -r -B -P $PIDFILE
	echo started
}

stop_service() {
	echo stop
	killall -15 $(cat "$PIDFILE") $$ rm -f "$PIDFILE"
	echo stoped
}

