# init.d/hostapd_action
#!/bin/sh /etc/rc.common
# OpenWrt init script for hostapd action hook
# Procd-managed, per-interface instances (no manual PIDs)

START=31
STOP=20
USE_PROCD=1

# --- Config ---
ACTION="/etc/hostapd_action"     # path to the action script
CONF_NS="hostapd_action"         # UCI config name
TAG="hostapd_action"             # syslog tag

# --- Libs ---
. /lib/functions.sh

log() { logger -t "$TAG" "$@"; }

# Globals populated from UCI
INCLUDE_REGEX=""
EXCLUDE_REGEX=""
IFACES=""

_append_iface() {
    # called by config_list_foreach network IFACE _append_iface
    local v="$1"
    [ -n "$v" ] && IFACES="$IFACES $v"
}

_load_config() {
    config_load "$CONF_NS"

    # Read network section
    config_get INCLUDE_REGEX network include ""
    config_get EXCLUDE_REGEX network exclude ""

    # list IFACE entries (optional)
    IFACES=""
    config_list_foreach network IFACE _append_iface
}

_discover_ifaces() {
    # Fallback discovery when no IFACE is provided in UCI
    # hostapd_cli interface prints:
    # "Available interfaces:"
    # wlan0
    # wlan1
    hostapd_cli interface 2>/dev/null | awk 'NR>1 && $1!="" {print $1}'
}

_filter_ifaces() {
    # Apply include/exclude regex if provided
    local out=""
    for i in $1; do
        [ -n "$INCLUDE_REGEX" ] && echo "$i" | grep -E -q "$INCLUDE_REGEX" || [ -z "$INCLUDE_REGEX" ] || continue
        [ -n "$EXCLUDE_REGEX" ] && echo "$i" | grep -E -q "$EXCLUDE_REGEX" && continue
        out="$out $i"
    done
    echo "$out"
}

_validate_runtime() {
    command -v hostapd_cli >/dev/null 2>&1 || {
        log "ERROR: hostapd_cli not found"
        return 1
    }
    [ -x "$ACTION" ] || {
        log "ERROR: action script missing or not executable: $ACTION"
        return 1
    }
    return 0
}

start_service() {
    _load_config

    # Build interface list: UCI IFACE list or discovery fallback
    local ifs="$IFACES"
    if [ -z "$ifs" ]; then
        ifs="$(_discover_ifaces)"
        [ -z "$ifs" ] && { log "WARNING: no interfaces discovered via hostapd_cli"; }
    fi

    # Apply include/exclude filters
    ifs="$(_filter_ifaces "$ifs")"

    _validate_runtime || return 1

    # Start one procd instance per interface, foreground (no -B)
    # procd supervises the process; hostapd_cli will reconnect (-r)
    local started=0
    for iface in $ifs; do
        [ -n "$iface" ] || continue
        procd_open_instance "ha_${iface}"
        procd_set_param command hostapd_cli -a "$ACTION" -r -i "$iface"
        procd_set_param respawn 5 60 0   # (threshold, timeout, retry) - reconnect/manage crashes
        procd_set_param stdout 1
        procd_set_param stderr 1
        procd_close_instance
        started=1
        log "Started action hook on interface: $iface"
    done

    if [ "$started" -ne 1 ]; then
        log "WARNING: nothing started (no matching interfaces). Check UCI list IFACE or adjust include/exclude."
    fi
}

stop_service() {
    # procd will terminate managed instances automatically
    log "Stopping all hostapd action instances"
}

reload_service() {
    # simple reload = restart with fresh config
    log "Reloading service (restart)"
    stop
    start
}

service_triggers() {
    # Reload on UCI config changes to this namespace
    procd_add_reload_trigger "$CONF_NS"
    # Also restart when hostapd restarts (optional, if available on your build)
    # procd_add_raw_trigger "config.change" 0 "/etc/init.d/hostapd" "restart"
}

